# supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[program:gunicorn]
command=/home/pejaab/sites/impudo-staging.ideacraft.de/virtualenv/bin/gunicorn impudo.wsgi:application --bind unix:/tmp/impudo-staging.ideacraft.de.socket --pid /home/pejaab/sites/impudo-staging.ideacraft.de/impudo/gunicorn.pid ;
directory=/home/pejaab/sites/impudo-staging.ideacraft.de/impudo/impudo-ui

[program:celeryd]
command=/home/pejaab/sites/impudo-staging.ideacraft.de/virtualenv/bin/celery worker -A impudo -l info
directory=/home/pejaab/sites/impudo-staging.ideacraft.de/impudo/impudo-ui
stdout_logfile=/var/log/celery/worker.log
stderr_logfile=/var/log/celery/worker.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
killasgroup=true
priority=998

[program:rabbitmq]
command=invoke-rc.d rabbitmq-server start
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
priority=999
